# ubuntu/etc/init/wlan-insmod.conf
# Wi-Fi
description			"Load the WLAN module"

# Triggers
start on started dbus
stop on stopping dbus

# Actions
exec insmod /system/lib/modules/wlan.ko
